# scripts/com.gmail.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gmail.sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/eddy/Developer/Langgraph/executive-assistant/.venv/bin/python3</string>
        <string>/Users/eddy/Developer/Langgraph/executive-assistant/scripts/gmail_daemon.py</string>
        <string>--watch</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>GWS_PROJECT</key>
        <string>gws-cli-20260428</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eddy/Library/Logs/gmail-sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eddy/Library/Logs/gmail-sync.err</string>
</dict>
</plist>
